# configs/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=usr/bin/mongod 
autorestart=true

[program:rabbitmq-server]
command=usr/sbin/rabbitmq-server 
autorestart=true